# cfg/installation/upstart/coffee.conf
description "upstart configuration for coffee"

start on started networking
stop on shutdown

respawn
setuid root
setgid root

exec coffee -cwb /var/apps/project/
